# system_backup_ver2/etc_systemd/audiophile-backend.service
[Unit]
Description=HQ Linux Music Player - FastAPI Backend
After=network.target mpd.service sound.target
Wants=mpd.service

[Service]
Type=simple
User=tysbox
Group=tysbox
WorkingDirectory=/home/tysbox/audiophile-web/backend
ExecStart=/home/tysbox/audiophile-web/backend/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
