# frontend/audiophile-frontend.service
[Unit]
Description=HQ Linux Music Player — Audiophile Frontend
After=network.target audiophile-backend.service
Wants=audiophile-backend.service
Requires=audiophile-backend.service

[Service]
Type=simple
User=tysbox
Group=tysbox
WorkingDirectory=/home/tysbox/HQ_Linux_Music_Player/frontend
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
Environment=PATH=/home/tysbox/.local/nodejs/node/bin:/usr/bin:/bin
ExecStartPre=/bin/sh -c 'mkdir -p /home/tysbox/HQ_Linux_Music_Player/frontend/.next/standalone/.next && rm -rf /home/tysbox/HQ_Linux_Music_Player/frontend/.next/standalone/.next/static && cp -r /home/tysbox/HQ_Linux_Music_Player/frontend/.next/static /home/tysbox/HQ_Linux_Music_Player/frontend/.next/standalone/.next/ && rm -rf /home/tysbox/HQ_Linux_Music_Player/frontend/.next/standalone/public && cp -r /home/tysbox/HQ_Linux_Music_Player/frontend/public /home/tysbox/HQ_Linux_Music_Player/frontend/.next/standalone/'
ExecStart=/home/tysbox/.local/nodejs/node/bin/node /home/tysbox/HQ_Linux_Music_Player/frontend/.next/standalone/server.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=audiophile-frontend

[Install]
WantedBy=multi-user.target