# deployment/systemd/audiophile-backend.service.d/override.conf
[Service]
# Ensure backend always sees the canonical MPD port
Environment=MPD_PORT=6600
