# deployment/systemd/audiophile-backend.service.d/override.conf
[Service]
# Pre-check MPD socket consistency before starting the service.
ExecStartPre=/home/tysbox/HQ_Linux_Music_Player/backend/scripts/check_mpd_socket.sh
# Ensure backend always sees the canonical MPD port
Environment=MPD_PORT=6600
