# config/systemd/fifo-drain.service
[Unit]
Description=MPD FIFO Drain - Prevents MPD freeze when no reader
After=mpd.service
Wants=mpd.service

[Service]
Type=simple
User=tysbox
Group=tysbox
ExecStart=/home/tysbox/HQ_Linux_Music_Player/backend/scripts/fifo_drain.sh
Restart=always
RestartSec=1
StandardOutput=null
StandardError=journal

[Install]
WantedBy=multi-user.target
