# config/systemd/loopback-drain.service
[Unit]
Description=ALSA Loopback drain — MPDフリーズ防止用ドレイン
# MPD起動後に開始。CamillaDSP または pure-passthrough が Loopback capture 側を
# 開いている間はそちらが優先され、このサービスは自動的に上書きされる。
After=sound.target mpd.service
Wants=mpd.service

[Service]
Type=simple
# Loopback capture 側 (hw:Loopback,1,0) を /dev/null に流し続ける。
# MPD が Loopback playback 側に書き込んでもバッファが詰まらないようにする。
ExecStart=/usr/bin/arecord \
    -D hw:Loopback,1,0 \
    -f S32_LE \
    -r 192000 \
    -c 2 \
    -t raw \
    --quiet \
    /dev/null
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
