# deployment/systemd/mpd.socket.d/override.conf
[Socket]
ListenStream=
ListenStream=/run/mpd/socket
ListenStream=127.0.0.1:6601